# cron.d/checker
*/5 * * * *    root    /usr/local/bin/checker/checker.sh 60
