# cron.d/checker
*/5 14-23 * * *    root    /usr/local/bin/checker/checker.sh
