# upstart/checker_off.conf
start on starting rc
task
exec /usr/local/bin/checker/log.sh off
